# build/iso/live-build/config/includes.chroot/etc/systemd/system/clawgress-live-selftest.service
[Unit]
Description=Clawgress Live ISO Self-Test
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/clawgress-live-selftest.sh
RemainAfterExit=true
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
